# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ---------------------------------------------------------------------------
; PHP-FPM — serves Laravel via Nginx
; ---------------------------------------------------------------------------
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/php-fpm.out.log
stderr_logfile=/var/log/php-fpm.err.log

; ---------------------------------------------------------------------------
; Nginx — reverse proxy & static file server
; ---------------------------------------------------------------------------
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/nginx.out.log
stderr_logfile=/var/log/nginx.err.log

; ---------------------------------------------------------------------------
; Laravel Queue Worker — processes background jobs (e.g. mail, notifications)
; ---------------------------------------------------------------------------
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2                ; Jalankan 2 worker secara paralel
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/queue-worker.out.log
stderr_logfile=/var/log/queue-worker.err.log
stopwaitsecs=60
